# scripts/adsb_stats.cron
*/30 * * * * root /usr/bin/bash /home/pi/stats.sh
